# services/pi-audio-cast-display.service
[Unit]
Description=Pi Audio Cast Display (minimal X, no desktop)
After=network.target sound.target
Wants=network.target
Before=display-manager.service lightdm.service
Conflicts=display-manager.service lightdm.service

[Service]
Type=simple
User=@RUN_USER@
Group=@RUN_USER@
SupplementaryGroups=video tty input
WorkingDirectory=@APP_DIR@
Environment=HOME=/home/@RUN_USER@
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-@APP_DIR@/.env
# Switch VT via chvt (systemd has no controlling tty). xinit runs as app user on vt7.
ExecStartPre=+/bin/chvt 7
ExecStart=/usr/bin/xinit @APP_DIR@/scripts/xinitrc -- :0 vt7 -keeptty -nocursor
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
